# Module1/Module1_Nginx/node-backend/node-backend-api.service
[Unit]
Description=Module 1 Nginx lab — Node demo API (Express)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/app
Environment=NODE_ENV=production
Environment=PORT=4000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
